# scripts/server/backup-flarum.cron
# backup flarum daily at midnight
@midnight /home/bitnami/devops/scripts/server/backup-flarum.sh
